# libs/avahi/files/avahi-daemon.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=61

BIN=avahi-daemon
DEFAULT=/etc/default/$BIN
OPTIONS="-D -f /rom/mnt/cust/etc/avahi/avahi-daemon.conf"
RUN_D=/var/run/$BIN

start() {
	[ -f $DEFAULT ] && . $DEFAULT
	mkdir -p $RUN_D
	/rom/mnt/cust/usr/sbin/$BIN $OPTIONS
}

stop() {
  $BIN -k
}

reload() {
  $BIN -r
}

